# telegram-memo-keeper/assets/systemd/memo-keeper.service
[Unit]
Description=MemoKeeper Telegram Bot
After=network.target redis.service

[Service]
Type=simple
User=memokeeper
WorkingDirectory=/opt/memo-keeper
Environment=PYTHONPATH=/opt/memo-keeper/scripts
EnvironmentFile=/opt/memo-keeper/.env
ExecStart=/opt/memo-keeper/venv/bin/python scripts/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
